# core/etc/profile.d/custom.sh
export ROUTER_USB_DIR=${USB_DIR}
export PATH=${USB_DIR}/mi_docker/docker-binaries:/data/usr/bin:${PATH}